# scripts/com.goldenhour.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.goldenhour.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/gregochr/IdeaProjects/goldenhour/scripts/backup-postgres.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gregochr/goldenhour-backups/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/gregochr/goldenhour-backups/launchd-stderr.log</string>
</dict>
</plist>
